# freebsd/etc/rc.local
zfs load-key storage
zfs mount storage/home
